# deploy/gcp/worker.service
[Unit]
Description=FX AutoBuy Background Worker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fxautobuy
EnvironmentFile=/opt/fxautobuy/.env
ExecStart=/opt/fxautobuy/venv/bin/python3 scripts/sim_runner.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
